# dist/ddnss-update.service
[Unit]
Description=ddnss automatic update

[Service]
# User=ddnss
User=root
Type=oneshot
ExecStart=/usr/bin/ddnss-update
